# mox/config/hypr/service/wallpaper.timer
[Unit]
Description=Wallpaper Service

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
